# projects/S812/packages/kodi/system.d/gui-switch.service
[Unit]
Description=GUI switch service
Before=kodi.service
After=var.mount userconfig.service

[Service]
Type=oneshot
ExecStart=/usr/bin/gui-switch
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
