# projects/S805/filesystem/usr/lib/systemd/system/zram.service
[Unit]
Description=ZRAM service

[Service]
Type=simple
ExecStart=/bin/sh -c 'exec /usr/bin/zram.start'

[Install]
WantedBy=swap.service
