# packages/xmedia/torrent/ptv/system.d/ptv.service
[Unit]
Description=Puzzle-TV IPTV channel aggregator
Requires=network-online.target
After=network-online.target
Before=tvheadend.service kodi.service
ConditionPathExists=/storage/.cache/services/ptv.conf

[Service]
ExecStart=/storage/.config/ptv3/server.py
ExecStop=/bin/kill -TERM $MAINPID
TimeoutStopSec=3
Restart=always
RestartSec=2
StartLimitInterval=0

[Install]
WantedBy=multi-user.target
